# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
#echo 2 > /sys/module/hid_apple/parameters/fnmode
echo 5 > /sys/class/backlight/acpi_video0/brightness
