# etc/systemd/system/mnt-readyshare.automount
[Unit]
Description=Automount Readyshare

[Automount]
Where=/mnt/readyshare

[Install]
WantedBy=multi-user.target
